# deploy/systemd/cogno-follower.service
# cogno-chain follower — prod-readiness Phase 1 (D1: read-only helper).
# Identity binding is now the permissionless on-chain self-proof `cognoGate.link_identity_signed`, so the
# follower NO LONGER verifies or writes bindings. This service is a small read-only helper: HTTP API on
# PORT (default 8090) exposing /health, /metrics, and /nonce (the bind-payload + live genesis). It holds
# NO signing key and forks no Node submitter — it only reads the node genesis via NODE_HTTP.

[Unit]
Description=cogno-chain follower (read-only identity helper)
Documentation=https://github.com/logical-mechanism/cogno
After=network-online.target cogno-node.service
Wants=network-online.target cogno-node.service

[Service]
Type=exec
User=cogno
Group=cogno
WorkingDirectory=/opt/cogno/services/cogno-follower
EnvironmentFile=/etc/cogno/cogno.env
ExecStart=/opt/cogno/services/cogno-follower/.venv/bin/python follower.py
Restart=always
RestartSec=10

StateDirectory=cogno
StateDirectoryMode=0700

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/cogno

[Install]
WantedBy=multi-user.target
