# deploy/systemd/cogno-relayer.service
# cogno-chain anchor relayer — prod-readiness Phase 1.
# Every ANCHOR_EVERY finalized blocks: write that block's finalized state-root to Cardano as tx
# metadata, then record anchor_ack on-chain (via the 3-of-5 committee by default). Funded Cardano hot
# wallet + crash-recovery cursor live in /var/lib/cogno (StateDirectory) — see services/_shared/paths.mjs.
#
# Restart=always + the relayer's SIGTERM graceful shutdown (finishes the current cycle, never mid-tx)
# make `systemctl restart` safe: a stop can't abandon a submitted-but-unpersisted paid Cardano tx.

[Unit]
Description=cogno-chain anchor relayer (finalized state-root -> Cardano metadata -> anchor_ack)
Documentation=https://github.com/logical-mechanism/cogno
After=network-online.target cogno-node.service
Wants=network-online.target cogno-node.service
# Also needs Ogmios + Kupo reachable. Add `After=`/`Wants=` for your cardano-node/ogmios/kupo units
# here (names vary by install) — the relayer fails closed and retries with backoff if they are down.

[Service]
Type=exec
User=cogno
Group=cogno
WorkingDirectory=/opt/cogno/services/anchor-relayer
EnvironmentFile=/etc/cogno/cogno.env
# Pin the data dir explicitly (the wallet + anchor cursor + single-instance lock live here).
Environment=COGNO_DATA_DIR=/var/lib/cogno
# /usr/local/bin/node MUST be a real node v22 (NOT the snap node, whose stdout goes to /dev/null) —
# symlink it in deploy/README.md "node binary".
ExecStart=/usr/local/bin/node relayer.mjs
Restart=always
RestartSec=10
# Give the graceful shutdown room to finish an in-flight Cardano confirmation before SIGKILL.
TimeoutStopSec=180

StateDirectory=cogno
StateDirectoryMode=0700

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/cogno

[Install]
WantedBy=multi-user.target
