# deploy/systemd/cogno-indexer.service
# cogno-chain SubQuery indexer (ingest) — prod-readiness Phase 1.
# Reads finalized blocks from the local archive node and writes the indexed schema into Postgres. Health
# on :3001. Requires Postgres reachable and the `cogno` role/db provisioned (see deploy/README.md).
#
# Restart=always handles the startup race: if Postgres or the node is not up yet the process exits and
# is restarted until they are.

[Unit]
Description=cogno-chain SubQuery indexer (ingest)
Documentation=https://github.com/logical-mechanism/cogno
After=network-online.target postgresql.service cogno-node.service
Wants=network-online.target cogno-node.service
Requires=postgresql.service

[Service]
Type=exec
User=cogno
Group=cogno
WorkingDirectory=/opt/cogno/services/indexer
EnvironmentFile=/etc/cogno/cogno.env
# @subql/node refuses to start unless TZ is UTC (deterministic timestamps).
Environment=TZ=UTC
ExecStart=/opt/cogno/services/indexer/node_modules/.bin/subql-node -f . --db-schema=cogno --batch-size=20 --workers=1 --port=3001
Restart=always
RestartSec=10

NoNewPrivileges=true
# `full` (not strict): subql tooling may write cache under its project dir in /opt; /etc stays read-only.
ProtectSystem=full
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
