# deploy/systemd/cogno-node.service
# cogno-chain validator node (Aura + GRANDPA) — prod-readiness Phase 1.
#
# A PERSISTENT, operator-keyed node — NOT `--dev` (which is an ephemeral //Alice single-authority dev
# chain). Session keys (aura sr25519 + grandpa ed25519) must be inserted into the keystore under
# <base-path>/chains/<id>/keystore BEFORE first start — see deploy/README.md ("validator keys").
# --chain points at the raw chainspec built by `cogno-chain-node gen-chainspec` (operator-owned genesis).
#
# Install: copy to /etc/systemd/system/, edit the literal paths/--name below to match your install,
#          then `systemctl daemon-reload && systemctl enable --now cogno-node`.

[Unit]
Description=cogno-chain node (Aura+GRANDPA validator)
Documentation=https://github.com/LogicalMechanism/cogno-chain/blob/main/deploy/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=cogno
Group=cogno
# The observer (in-protocol Cardano weight, D4) reads the read-only Cardano db-sync at DBSYNC_URL from the
# shared EnvironmentFile. The leading `-` makes it OPTIONAL: a missing file or an unset DBSYNC_URL is NOT
# fatal — the observer simply ABSTAINS (empty observation, non-fatal) and the node keeps producing +
# finalizing. The node consumes only DBSYNC_URL (+ optional RUST_LOG) from this file; db-sync must be FULL /
# tx_in-enabled (see deploy/README.md + .env.example). Without this line the observer logs
# "no DBSYNC_URL/DBSYNC set — abstaining" and credits nothing.
EnvironmentFile=-/etc/cogno/cogno.env
# --state-pruning/--blocks-pruning archive: the node serves ALL reads (feed / thread / search / profile)
#   from its own runtime API, so it must retain history; and with a single validator the node DB is the
#   sole copy of chain history (back it up — see deploy/README). A non-archive node can't answer them.
# RPC/Prometheus bind to localhost by default; the frontend reads over 127.0.0.1. `--rpc-methods safe`
#   is ON below: even on a localhost bind, `auto` would still expose UNSAFE methods (author_insertKey,
#   author_rotateKeys, peer manipulation) to any same-host browser page, and `--rpc-cors all` drops the
#   origin guard — the static frontend needs only SAFE methods (state_*/chain_*/system_health/
#   author_submitExtrinsic). Narrow `--rpc-cors` to your frontend origin(s) where you can, and put any
#   public exposure (browser frontend, remote metrics) behind a TLS reverse proxy.
# --force-authoring: REQUIRED for the single-validator default topology — a lone validator has no
#   peers to confirm against and will not author without it. REMOVE it for a multi-validator network
#   (authoring should begin only once peers connect).
ExecStart=/usr/local/bin/cogno-chain-node run \
  --validator \
  --force-authoring \
  --name %H \
  --chain /etc/cogno/chainspec.raw.json \
  --base-path /var/lib/cogno/node \
  --state-pruning archive \
  --blocks-pruning archive \
  --port 30333 \
  --rpc-port 9944 \
  --rpc-methods safe \
  --rpc-cors all \
  --prometheus-port 9615 \
  --no-mdns
Restart=always
RestartSec=5
LimitNOFILE=65536
TimeoutStopSec=120

# StateDirectory creates/owns /var/lib/cogno (0700, cogno:cogno); the node DB lives in node/ under it.
StateDirectory=cogno
StateDirectoryMode=0700

# Sandboxing — the node only needs to write its DB under /var/lib/cogno.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/lib/cogno
# Kernel / privilege-surface hardening. NOTE: deliberately NO MemoryDenyWriteExecute — the runtime is
# executed by a JIT (wasmtime) that needs W^X pages, and denying W+X would break block execution.
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
# Further low-cost surface reduction: a validator binds only high ports and needs no privileges or
# devices. SystemCallFilter uses the @system-service baseline (wasmtime's W^X JIT works under it); if a
# future release trips it, relax with `systemd-analyze security cogno-node` and re-check. Deny @privileged
# only (NOT @resources — the async runtime may touch scheduler/affinity calls).
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK
CapabilityBoundingSet=
AmbientCapabilities=
PrivateDevices=true
ProtectHostname=true
ProtectClock=true
ProtectProc=invisible
ProcSubset=pid
RemoveIPC=true
UMask=0077

[Install]
WantedBy=multi-user.target
